# profile.d/002_histsize.sh
#!/bin/sh
#
#	Configuration of history file size
#	~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~
#
#	@(#)	[MB] de_002_histsize.sh	Version 1.1 du 19/11/06 - 
#

export HISTSIZE=2000
